# systemd/dotaml-live-retrain.timer
[Unit]
Description=Nightly dotaml-live retrain cycle

[Timer]
# Mon-Sat at 04:00; Sunday is the from-scratch weekly retrain instead.
OnCalendar=Mon..Sat *-*-* 04:00:00
# Catch up if the machine was off at the scheduled time.
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
